# supervisord.conf
[supervisord]
nodaemon=true

[program:spacewalk-setup]
command=spacewalk-setup --answer-file /tmp/answerfile.txt
autorestart=false
